# systemd/intel-speech-to-text.service
[Unit]
Description=Intel Speech-to-Text Key Listener
After=ydotoold.service
Wants=ydotoold.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/guido/intel-speech-to-text
ExecStart=/home/guido/anaconda3/envs/intel-stt/bin/python /home/guido/intel-speech-to-text/key_listener.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
